# finance-bot.service
[Unit]
Description=Finance Telegram Bot
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/finance-bot
ExecStart=/usr/bin/python3 /root/finance-bot/bot.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
